# templates/letsencrypt.timer.j2
[Unit]
Description=Daily Letsencrypt certificate renewal

[Timer]
OnCalendar=daily
