# templates/letsencrypt.service.j2
[Unit]
Description=Letsencrypt certificate renewal
After=network.service

[Service]
Type=oneshot
ExecStart=/usr/bin/letsencrypt renew
