# files/observium-poller.service
[Unit]
Description=Observium Poller

[Service]
Type=oneshot
ExecStart=/opt/observium/poller-wrapper.py 10
