# files/observium-housekeeping.timer
[Unit]
Description=Daily Observium housekeeping

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
