# files/observium-housekeeping.service
[Unit]
Description=Observium housekeeping

[Service]
Type=oneshot
ExecStart=/opt/observium/housekeeping.php -ay
