# files/observium-housekeeping.timer
[Unit]
Description = Run Observium housekeeping

[Timer]
OnCalendar = daily

[Install]
WantedBy = timers.target
